# systemd/abfahrt-watchdog.service
[Unit]
Description=Abfahrt health check

[Service]
Type=oneshot
ExecStart=/home/pi/abfahrt/systemd/healthcheck.sh
